# app/crontab
 * 1 * * * python3 /code/app/scrapping.py >> /code/app/log_cron 2>&1
